# templates/profile.d/codeserver.sh
if [ -n "${BASH_VERSION:-}" ] && [ -f "$HOME/.scripts/index.sh" ]; then
  . "$HOME/.scripts/index.sh"
fi
